# data/data/com/termux/files/usr/var/service/XOXO/run
cat > "$PREFIX/var/service/xoxo/run" << 'EOF'
#!/data/data/com.termux/files/usr/bin/sh

# Activar entorno de Python si existe
if [ -f "$PREFIX/bin/python" ]; then
    PY=python
else
    PY=python3
fi

LOGFILE="$HOME/xoxo/xoxo.log"

echo "Iniciando XOXO Python..." >> "$LOGFILE"

# Bucle para asegurar que si Python falla, el servicio reinicie
while true; do
    $PY $HOME/xoxo/xoxo.py
    echo "XOXO Python se reinició tras un fallo..." >> "$LOGFILE"
    sleep 2
done
EOF
